# 0.X.X Personal MVP/0.5.1/iso-build/overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin noctem --noclear %I $TERM
