# outdated/parent/systemd/noctem-babysit.service
[Unit]
Description=Noctem Babysitting Report
Documentation=https://github.com/noctem/noctem
After=network.target

[Service]
Type=oneshot
# Update this path to your Noctem installation
ExecStart=/usr/bin/python3 /path/to/noctem/parent/scheduler.py --report --analyze
# Run as your user
User=%i
# Working directory
WorkingDirectory=/path/to/noctem

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
